# files/etc/init.d/dns-switchy
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

USER_CONFIG=/etc/dns-switchy/config.yaml
RUNTIME_CONFIG=/var/run/dns-switchy.yaml

generate_runtime_config() {
	local listen_addr="$1"
	local http_port="$2"

	mkdir -p /var/run

	# UCI fields override yaml top-level addr:/http: lines
	{
		printf 'addr: "%s"\n' "${listen_addr}"
		printf 'http: ":%s"\n' "${http_port}"
		if [ -f "${USER_CONFIG}" ]; then
			sed -E '/^addr:[[:space:]]/d; /^http:[[:space:]]/d' "${USER_CONFIG}"
		else
			echo 'ttl: 5m'
			echo 'resolvers:'
			echo '  - name: default'
			echo '    type: forward'
			echo '    url: 114.114.114.114'
		fi
	} > "${RUNTIME_CONFIG}"
}

start_service() {
	. /lib/functions.sh
	config_load 'dns-switchy'

	local enabled listen_addr http_port
	config_get_bool enabled 'main' 'enabled' 1
	config_get listen_addr 'main' 'listen_addr' ':1053'
	config_get http_port 'main' 'http_port' '8080'

	[ "${enabled}" = '1' ] || return 1

	generate_runtime_config "${listen_addr}" "${http_port}"

	procd_open_instance dns-switchy
	procd_set_param command /usr/bin/dns-switchy -c "${RUNTIME_CONFIG}"
	procd_set_param respawn "${respawn_threshold:-3600}" "${respawn_timeout:-5}" "${respawn_retry:-5}"
	procd_set_param limits core="unlimited"
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param user root
	procd_set_param pidfile /var/run/dns-switchy.pid
	procd_set_param file "${USER_CONFIG}" /etc/config/dns-switchy
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger 'dns-switchy'
}
