# dist/com.test4x.dns-switchy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.test4x.dns-switchy</string>
        <key>WorkingDirectory</key>
        <string>/Users/{user}/bin/dns</string>
        <key>ProgramArguments</key>
        <array>
		    <string>/Users/{user}/bin/dns/dns-switchy</string>
		    <string>-c</string>
		    <string>config.yaml</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
        <key>ThrottleInterval</key>
        <integer>5</integer>
        <key>StandardErrorPath</key>
	    <string>/usr/local/var/log/dns-switchy.log</string>
	    <key>StandardOutPath</key>
	    <string>/usr/local/var/log/dns-switchy.log</string>
    </dict>
</plist>